# pkg-build-monolithic/scripts/postinstall
#!/bin/bash
# Monolithic Package Postinstall - everything pre-compiled

set +e

LOG_FILE="/var/log/speedmonitor-install.log"
INSTALL_DIR="/usr/local/speedmonitor"

log() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') - $1" | tee -a "$LOG_FILE"
}

log "=== Speed Monitor Monolithic Postinstall ==="

# Load configuration
if [ -f "${INSTALL_DIR}/lib/config.sh" ]; then
    source "${INSTALL_DIR}/lib/config.sh"
    log "Configuration loaded"
fi

# Set up for all users
log "Setting up Speed Monitor for all users..."

for user_home in /Users/*; do
    if [ ! -d "$user_home" ]; then
        continue
    fi

    username=$(basename "$user_home")

    # Skip system users
    if [ "$username" = "Shared" ] || [ "$username" = "Guest" ]; then
        continue
    fi

    log "Configuring for user: $username"

    # Create directories
    sudo -u "$username" mkdir -p "$user_home/${USER_DATA_DIR}"
    sudo -u "$username" mkdir -p "$user_home/${USER_CONFIG_DIR}"
    sudo -u "$username" mkdir -p "$user_home/${USER_BIN_DIR}"
    sudo -u "$username" mkdir -p "$user_home/Library/LaunchAgents"

    # Link binaries to user bin
    sudo -u "$username" ln -sf "${INSTALL_DIR}/bin/speed_monitor.sh" "$user_home/${USER_BIN_DIR}/speed_monitor.sh"
    sudo -u "$username" ln -sf "${INSTALL_DIR}/bin/wifi_info" "$user_home/${USER_BIN_DIR}/wifi_info"
    sudo -u "$username" ln -sf "${INSTALL_DIR}/bin/speedtest-cli" "$user_home/${USER_BIN_DIR}/speedtest-cli"

    # Generate device ID if not exists
    if [ ! -f "$user_home/${USER_CONFIG_DIR}/device_id" ]; then
        hw_uuid=$(ioreg -rd1 -c IOPlatformExpertDevice | awk '/IOPlatformUUID/ { print $3 }' | tr -d '"')
        device_id=$(echo "$hw_uuid-$username" | shasum -a 256 | cut -c1-16)
        echo "$device_id" | sudo -u "$username" tee "$user_home/${USER_CONFIG_DIR}/device_id" > /dev/null
        log "Generated device ID for $username: $device_id"
    fi

    # Create LaunchAgent plist
    cat > "$user_home/Library/LaunchAgents/com.speedmonitor.plist" << PLIST_EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.speedmonitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>$user_home/${USER_BIN_DIR}/speed_monitor.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>$user_home/${USER_DATA_DIR}/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>$user_home/${USER_DATA_DIR}/launchd_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>$user_home/${USER_BIN_DIR}:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SPEED_MONITOR_SERVER</key>
        <string>${SERVER_URL}</string>
    </dict>
</dict>
</plist>
PLIST_EOF

    chown "$username" "$user_home/Library/LaunchAgents/com.speedmonitor.plist"
    log "Created LaunchAgent for $username"

    # Load LaunchAgent (only if user is logged in)
    if pgrep -u "$username" &> /dev/null; then
        sudo -u "$username" launchctl load "$user_home/Library/LaunchAgents/com.speedmonitor.plist" 2>/dev/null || true
        log "Loaded LaunchAgent for $username"
    else
        log "User $username not logged in - LaunchAgent will load on next login"
    fi

    # SpeedMonitor.app is installed directly to /Applications by the package
    log "SpeedMonitor.app available at /Applications/SpeedMonitor.app"
done

# Run initial speed test for console user
CURRENT_USER=$(stat -f "%Su" /dev/console)
if [ -f "/Users/${CURRENT_USER}/${USER_BIN_DIR}/speed_monitor.sh" ]; then
    log "Running initial speed test for $CURRENT_USER..."
    sudo -u "$CURRENT_USER" "/Users/${CURRENT_USER}/${USER_BIN_DIR}/speed_monitor.sh" &
fi

log "=== Postinstall Complete ==="
log "Speed Monitor installed successfully!"
log "Dashboard: ${SERVER_URL}"

exit 0
